# selene_project.conf
[program:selene_project-app]
command=/path/to/selene_project_start
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/path/to/selene_project/logs/uvicorn-info.log